# com.predator.switch_model_on_idle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.predator.switch_model_on_idle</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>/Users/Shared/Predator_60/scripts/switch_model_on_idle.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/Shared/Predator_60/logs/switch_model_on_idle.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/Shared/Predator_60/logs/switch_model_on_idle_error.log</string>
</dict>
</plist>